# xeyes.cron
# Run xeyes every minute
*/1 * * * * root DISPLAY=:0.0 /usr/bin/xeyes &
